# root/etc/crontabs/root
0 * * * * /bin/bash /app/musicbrainz/admin/cron/slave.sh 
