# root/etc/services.d/redis/run
#!/usr/bin/with-contenv bash

umask 022

s6-setuidgid abc redis-server /config/redis.conf


