# services/postgres/run
#!/bin/bash
exec /sbin/setuser abc postgres -D $PGDATA -c config_file=$PGCONF/postgresql.conf

