# root/etc/services.d/server/run
#!/usr/bin/with-contenv bash

cd /app/musicbrainz || exit

s6-setuidgid abc start_server --port=5000 -- plackup -I lib -s Starlet \
--max-workers "${MAX_WORKERS}" > /dev/null 2>&1
